# files/etc/uci-defaults/99-change-default-shell-to-bash
#!/bin/sh
[ -f /bin/bash ] && chsh -s /bin/bash
