# files/etc/uci-defaults/99-remove-forced-zsh
#!/bin/bash
# Remove "/usr/bin/zsh" placed at the last line of /etc/profile. Who f**king places the shit?
sed -i.bak '/^\/usr\/bin\/zsh/d' /etc/profile
