# files/etc/uci-defaults/99-tmux-switch-client-sh-chmod-755
#!/bin/sh
# Auto add execute permission to /etc/tmux.switch-client.sh (tmux.conf calls this without "bash" prefix)
[ -f "/etc/tmux.switch-client.sh" ] && chmod 755 "/etc/tmux.switch-client.sh"
